# systemd/sherpad.service
[Unit]
Description=Sherpa Daemon - VM, Container, and Unikernel Management Server
Documentation=https://github.com/bwks/sherpa
After=network-online.target docker.service libvirtd.service
Wants=network-online.target
Requires=docker.service libvirtd.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=sherpa
Group=sherpa
WorkingDirectory=/opt/sherpa

# Run in foreground so systemd manages the process lifecycle
ExecStart=/opt/sherpa/bin/sherpad start --foreground

# Restart policy
Restart=on-failure
RestartSec=5

# Environment configuration (optional file)
EnvironmentFile=-/opt/sherpa/env/sherpa.env

# Security hardening (basic level)
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
